# deploy/enp.supervisor.conf
[program:enp]
command=/bin/bash -c 'set -a; source /opt/enp/.env; set +a; exec /usr/bin/node dist/index.cjs'
directory=/opt/enp
user=root
autostart=true
autorestart=true
startretries=5
startsecs=5
environment=NODE_ENV="production",PORT="8003"
stderr_logfile=/var/log/enp.err.log
stdout_logfile=/var/log/enp.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
